# roles/fstrim/templates/fstrim.timer.j2
# fstrim.timer override.conf

[Unit]
Description=Discard unused filesystem blocks once a week

[Timer]
OnCalendar=weekly
