# roles/pacman/templates/paccache.timer.j2
# paccache.timer override.conf

[Unit]
Description=Discard unused packages weekly

[Timer]
OnCalendar=weekly
