# roles/pacman/templates/reflector.timer.j2
# reflector.timer override.conf

[Unit]
Description=Refresh Pacman mirrorlist weekly with Reflector

[Timer]
OnCalendar=weekly

